# supervisord.conf

# supervisord.conf
[supervisord]
nodaemon=true

[program:asterisk]
command=/usr/sbin/asterisk -f -U root

[program:flask]
directory=/app
command=python3 main.py
autostart=true
autorestart=true